# containers/services/hello-world/run
#!/bin/sh -e
# Show error to stdout
exec 2>&1

# Add shell command here
exec echo "Hello World!"
